# postflow-watchdog.service
[Unit]
Description=Watchdog for postflow_django cron health
After=docker.service
Requires=docker.service

[Service]
ExecStart=/bin/bash /home/ubuntu/postflow/watchdog.sh
StandardOutput=append:/home/ubuntu/logs/watchdog.log
StandardError=append:/home/ubuntu/logs/watchdog.log
Restart=always
RestartSec=5
User=ubuntu

[Install]
WantedBy=multi-user.target
